# maplestory-timekeeper.service.sample
[Unit]
Description=MapleStory Timekeeper

[Service]
User=
Group=
ExecStart=
WorkingDirectory=
RemainAfterExit=no
Restart=on-failure
RestartSec=2s

[Install]
WantedBy=multi-user.target